# packaging/systemd/nia-todo-backup.timer
[Unit]
Description=Daily nia-todo backup timer

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
Unit=nia-todo-backup.service

[Install]
WantedBy=timers.target
